# templates/agora-api/agora-api.conf
[program:agora-api]
command=/home/agoraapi/go/bin/agora-api
directory=/home/agoraapi/agora-api/
autostart=true
autorestart=true
user=agoraapi
redirect_stderr=true
stopasgroup=true
