# templates/eorchestra/election-orchestra.conf
[program:eorchestra]
command=/home/eorchestra/launch.sh
directory=/home/eorchestra/election-orchestra/
autostart=true
autorestart=true
user=eorchestra
redirect_stderr=true
stopasgroup=true