# authapi/templates/supervisor_authapi_celery_beat.conf
[program:authapi_celery_beat]
command=/home/authapi/authapi_celery_beat.sh
process_name=%(program_name)s
directory=/home/authapi/authapi/authapi/
autostart=true
autorestart=true
stopwaitsecs=3
user=authapi
environment=PATH="/home/authapi/env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
redirect_stderr=true
